# roles/xboxdrv/files/xboxdrv-daemon.service
[Unit]
Description=xboxdrv

[Service]
ExecStart=/usr/bin/xboxdrv --silent -D --dbus disabled --detach

[Install]
WantedBy=multi-user.target
